# deploy/gunicorn.diving.service
# systemd unit for Gunicorn (DivinG_BAZ) — tailor paths if needed
[Unit]
Description=gunicorn daemon for DivinG_BAZ
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/DivinG_BAZ
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/.local/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 app:app
Restart=always

[Install]
WantedBy=multi-user.target
