# deploy/gunicorn.service
# systemd service file for Gunicorn (DivinG_BAZ) on EC2/Linux
# Installation:
#   1. Edit paths below if needed
#   2. Copy to: sudo cp gunicorn.service /etc/systemd/system/
#   3. Enable: sudo systemctl daemon-reload && sudo systemctl enable gunicorn
#   4. Start: sudo systemctl start gunicorn
#   5. Check: sudo systemctl status gunicorn

[Unit]
Description=gunicorn daemon for Diving
After=network.target

[Service]
Type=notify
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/diving
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/.local/bin/gunicorn --workers 3 --bind 127.0.0.1:3000 --timeout 90 --access-logfile - --error-logfile - app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
